# pixie.service
[Unit]
Description=Pixie LED Matrix Controller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/pixie
ExecStart=/usr/bin/python3 -u run_pixie.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
